# packaging/debian/openvas-server.prerm
#!/bin/sh

set -e

if [ "$1" != "upgrade" ]; then
	if [ -x "/etc/init.d/openvas-server" ] ; then
		        if which invoke-rc.d >/dev/null 2>&1; then
		                invoke-rc.d openvas-server stop
			else
	        		/etc/init.d/openvas-server stop
			fi
	else
                start-stop-daemon --stop --quiet --oknodo --exec /usr/sbin/openvasd
	fi
fi


#DEBHELPER#


exit 0
