# .docker/e2e/gvmd/service/gvmd/run
#!/bin/sh

exec 2>&1
exec su -s /bin/sh gvmd -c "/usr/local/bin/start-gvmd"
