# .docker/e2e/gvmd/service/postgresql/run
#!/bin/sh

exec 2>&1
exec su -s /bin/sh postgres -c "POSTGRES_VERSION=15 /usr/local/bin/start-postgresql"
